# scripts/launchd/analyzer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.solana-copilot.analyzer</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update COPILOT_DIR to the actual project path on your Mac mini -->
    <string>/Users/YOUR_USERNAME/solana-copilot/.venv/bin/uvicorn</string>
    <string>src.services.analyzer_server:app</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8000</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/solana-copilot</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/solana-copilot/logs/analyzer.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/solana-copilot/logs/analyzer.err</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
